# docker/supervisord/php-fpm.ini
[program:php-fpm]
command = php-fpm -F
autorestart = true
stderr_logfile = /tmp/php.err
stdout_logfile = /tmp/php.log